# backup/rclone-backup.timer
[Unit]
Description=Ежедневный триггер CephFS-бэкапа (rclone-backup.service)
# Скрипт сам держит flock (/var/lock/backup.lock). Реальный цикл занимает
# 5-8 суток, поэтому большинство ежедневных срабатываний застанут предыдущий
# запуск ещё выполняющимся и корректно завершатся сразу с "уже выполняется"
# (это ожидаемо, не является сбоем). Новый цикл фактически стартует, как
# только предыдущий завершится и наступит ближайшее время срабатывания.

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
